# resources/supervisor/programs/cron.conf
[program:cron]
command=/usr/sbin/cron -f -L 15 ; run in forground and log everything
autostart=true
autorestart=true
redirect_stderr=true 
stdout_logfile=/var/log/supervisor/%(program_name)s.log ; log logs into file
startretries=5   ; max # of serial start failures (default 3)