# resources/supervisor/programs/nginx.conf
[program:nginx]
command=/usr/local/openresty/nginx/sbin/nginx -c /etc/nginx/nginx.conf -g "daemon off;"
priority=50
redirect_stderr=true 
;stdout_logfile=/var/log/supervisor/%(program_name)s.log ; log logs into file
stdout_logfile=/proc/1/fd/1 ; log error into stdout to see in docker logs
stdout_logfile_maxbytes=0
autostart=true ; start at supervisord start (default: true)
autorestart=true ; whether/when to restart (default: unexpected) 
startretries=5   ; max # of serial start failures (default 3)