# resources/supervisor/programs/novnc.conf
[program:novnc]
command=python -m websockify --web  %(ENV_RESOURCES_PATH)s/novnc/ 6901 localhost:5901
; the launch script cannot be terminated -> %(ENV_RESOURCES_PATH)s/novnc/utils/launch.sh --vnc localhost:5901 --listen 6901
directory=%(ENV_RESOURCES_PATH)s/novnc/utils/websockify/
autostart=true
autorestart=true
redirect_stderr=true 
stdout_logfile=/var/log/supervisor/%(program_name)s.log ; log logs into file
startretries=5   ; max # of serial start failures (default 3)