# resources/supervisor/programs/vncserver.conf
[program:vncserver]
command=%(ENV_RESOURCES_PATH)s/scripts/start-vnc-server.sh
priority=55
autostart=true
autorestart=true
redirect_stderr=true 
stdout_logfile=/var/log/supervisor/%(program_name)s.log ; log logs into file
startretries=5   ; max # of serial start failures (default 3)