# resources/supervisor/programs/vscode.conf
[program:vscode]
command=/usr/local/bin/code-server --port=8054 --disable-telemetry --user-data-dir=%(ENV_HOME)s/.config/Code/ --extensions-dir=%(ENV_HOME)s/.vscode/extensions/ --auth=none %(ENV_WORKSPACE_HOME)s/
autostart=true
autorestart=true
redirect_stderr=true 
stdout_logfile=/var/log/supervisor/%(program_name)s.log ; log logs into file
startretries=5   ; max # of serial start failures (default 3)