# spark-flavor/resources/zeppelin-service.conf
[program:zeppelin]
command=%(ENV_RESOURCES_PATH)s/tools/zeppelin.sh --port=8072
redirect_stderr=true 
stdout_logfile=/var/log/supervisor/%(program_name)s.log ; log logs into file
autostart=true ; start at supervisord start (default: true)
autorestart=true ; whether/when to restart (default: unexpected) 
startretries=5   ; max

